# systemd/gpustat-dashboard.service
[Unit]
Description=GPUStat Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=shimdx
Group=shimdx
WorkingDirectory=/home/shimdx/Workspace/gpustat
Environment=HOME=/home/shimdx
Environment=NODE_ENV=production
ExecStart=/home/shimdx/Workspace/gpustat/scripts/start-prod.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
